# pyalertec2.service
[Unit]
Description=PyAlertEC2 System Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/pyalertec2/monitor --webhook YOUR_SLACK_WEBHOOK_URL
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
